# all_applications_runner/service_files/app-nimble_books-imprint_builder.service
[Unit]
Description=Imprint Builder - Nimble Books LLC
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/home/wfz/nimble/repos/imprint-builder
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/wfz/nimble/repos/imprint-builder
ExecStart=/usr/bin/python3 -m streamlit run  app.py --server.port=8507
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target