# systemd_services/daily-engine.service
[Unit]
Description=Daily Engine - Personal Time Management
After=network.target

[Service]
Type=simple
User=wfzimmerman
Group=wfzimmerman
WorkingDirectory=/home/wfzimmerman/xcu_my_apps/personal
Environment="PATH=/home/wfzimmerman/xcu_my_apps/personal/.venv/bin:/home/wfzimmerman/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/wfzimmerman/xcu_my_apps"
ExecStart=/usr/local/bin/uv run streamlit run daily_engine.py --server.port=8509 --server.address=127.0.0.1 --server.headless=true
Restart=always
RestartSec=10
StandardOutput=append:/home/wfzimmerman/xcu_my_apps/personal/logs/daily-engine.log
StandardError=append:/home/wfzimmerman/xcu_my_apps/personal/logs/daily-engine-error.log

[Install]
WantedBy=multi-user.target
